# services/suppress-gpe4E.service
[Unit]
Description=Diables GPE4E, which causes a ACPI storm on MacBooks and prevents low power states on CPU

[Service]
ExecStart=/usr/bin/bash -c 'echo "disable" > /sys/firmware/acpi/interrupts/gpe4E'

[Install]
WantedBy=multi-user.target
